# scripts/metrics.service
# Install:
#   sudo cp metrics.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable metrics
#   sudo systemctl start metrics
#
# Check status:
#   sudo systemctl status metrics
#   journalctl -u metrics -f

[Unit]
Description=DOTM Metrics Server
After=network.target

[Service]
Type=simple
WorkingDirectory=/root/dotm-sniper
ExecStart=/usr/bin/python3 src/metrics_server.py 8765
Restart=on-failure
RestartSec=30
StandardOutput=journal
StandardError=journal
Environment=LOG_FORMAT=json
EnvironmentFile=/root/dotm-sniper/.env

[Install]
WantedBy=multi-user.target
